# files/debian/supervisor/storage-service.conf
[program:storage-service]
command = uwsgi_python --ini=/etc/uwsgi/apps-enabled/storage.ini 
